# systemd/mission-control.service
[Unit]
Description=OpenClaw Mission Control
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/.openclaw/workspace/mission-control/backend
Environment=MC_HOST=127.0.0.1
Environment=MC_PORT=8787
Environment=SPARK_OLLAMA=http://127.0.0.1:11434
Environment=SPARK_METRICS_URL=http://127.0.0.1:8766/metrics
ExecStart=/home/openclaw/.openclaw/workspace/mission-control/.venv/bin/python /home/openclaw/.openclaw/workspace/mission-control/backend/app.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
